# templates/cultivator.timer.erb
[Unit]
Description=Schedule cultivator runs

[Timer]
OnBootSec=<%= @bootdelay %>
OnUnitActiveSec=<%= @frequency %>

[Install]
WantedBy=multi-user.target
